# ircbot.service
[Unit]
Description=Simple IRC bot
Documentation=https://github.com/cbosoft/ircbot

[Service]
ExecStart=ircbot
